# cron/crontab
*/2 * * * * /home/marcel/.dotfiles/cron/hibernate-on-low-battery.sh
@hourly /home/marcel/.dotfiles/cron/journal-reminder.sh
