# systemd/blueman-applet.service
[Unit]
Description=a tray applet for managing bluetooth
After=graphical.target
OnFailure=notify@%i.service

[Service]
ExecStart=/usr/bin/blueman-applet

Restart=on-failure
RestartSec=30s
StartLimitBurst=3

[Install]
WantedBy=default.target
